# config/nexuserofs-snapshotter.service
[Unit]
Description=Nexus EROFS Snapshotter for containerd
Documentation=https://github.com/aledbf/nexus-erofs
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nexus-erofs-snapshotter \
    --address /run/nexus-erofs-snapshotter/snapshotter.sock \
    --root /var/lib/nexus-erofs-snapshotter \
    --containerd-address /run/containerd/containerd.sock \
    --log-level info
Restart=always
RestartSec=5
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity
TasksMax=infinity
OOMScoreAdjust=-999

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
